# luci-app-wizard/root/etc/uci-defaults/40_luci-app-wizard
#!/bin/sh

[ ! -f "/usr/share/ucitrack/luci-app-wizard.json" ] && {
	cat > /usr/share/ucitrack/luci-app-wizard.json << EOF
		{
		    "config": "wizard",
		    "init": "wizard"
		}
	EOF
}

[ -f "/usr/share/luci/menu.d/luci-app-wizard.json" ] || {
	cat > /usr/share/luci/menu.d/luci-app-wizard.json << EOF
		{
		    "admin/system/wizard": {
		        "title": "Wizard",
		        "order": 10,
		        "action": {
		            "type": "firstchild",
		        },
		        "depends": {
		            "acl": [ "luci-app-wizard" ],
		            "uci": { "wizard": true }
		        }
		    },

		    "admin/system/wizard/wizard": {
		        "title": "Configuration",
		        "order": 10,
		        "action": {
		            "type": "view",
		            "path": "wizard/wizard"
		        }
		    },

		    "admin/system/wizard/file": {
		        "title": "Files",
		        "order": 20,
		        "action": {
		            "type": "view",
		            "path": "wizard/file"
		        }
		    }
		}
	EOF
}

[ ! -f "/usr/share/rpcd/acl.d/luci-app-wizard.json" ] && {
	cat > /usr/share/rpcd/acl.d/luci-app-wizard.json << EOF
		{
		    "luci-app-wizard": {
		        "description": "Grant UCI access for luci-app-wizard",
		        "read": {
		            "uci": [ "wizard" ],
		            "file": {
		                "/etc/init.d/wizard reconfig": [ "exec" ],
		                "/etc/hosts": [ "read" ],
		                "/etc/rc.local": [ "read" ],
		                "/etc/config/dhcp": [ "read" ],
		                "/etc/dnsmasq.conf": [ "read" ],
		                "/etc/config/uhttpd": [ "read" ],
		                "/etc/config/network": [ "read" ],
		                "/etc/config/firewall": [ "read" ],
		                "/etc/config/wireless": [ "read" ]
		            },
		            "ubus": {
		                "file": [ "read" ],
		                "rc": [ "list" ]
		            }
		        },
		        "write": {
		            "file": {
		                "/etc/hosts": ["write"],
		                "/etc/rc.local": ["write"],
		                "/etc/config/dhcp": ["write"],
		                "/etc/dnsmasq.conf": ["write"],
		                "/etc/config/uhttpd": ["write"],
		                "/etc/config/network": ["write"],
		                "/etc/config/firewall": ["write"],
		                "/etc/config/wireless": ["write"],
		                "/etc/init.d/network reload": ["exec"],
		                "/etc/init.d/cron reload": ["exec"],
		                "/etc/init.d/firewall reload": ["exec"],
		                "/etc/init.d/dnsmasq reload": ["exec"],
		                "/etc/init.d/uhttpd reload": ["exec"]
		            },
		            "ubus": {
		                "file": ["write"],
		                "rc": ["init"]
		            }
		        }
		    }
		}
	EOF
}

touch /etc/config/wizard
uci -q batch <<-EOF >/dev/null
	delete ucitrack.@wizard[-1]
	add ucitrack wizard
	set ucitrack.@wizard[-1].init=wizard
	commit ucitrack
	set wizard.default=wizard
	commit wizard
	set upnpd.config.enabled=1
	commit upnpd
	set system.@system[0].hostname=OpenWrt
	commit system
	set luci.main.mediaurlbase=/luci-static/bootstrap
	commit luci
EOF
chmod +x /etc/init.d/wizard 2>/dev/null
rm -f /tmp/luci-*
for file in dhcp network firewall wireless; do
    [ ! -f /etc/config/bk_$file -a -s /etc/config/$file ] && cp /etc/config/$file /etc/config/bk_$file
done

sed -i 's/root:.*/root:$1$pn1ABFaI$vt5cmIjlr6M7Z79Eds2lV0:16821:0:99999:7:::/g' /etc/shadow
exit 0
